# file_roots/mint_18-2/network/etc/network/if-up.d/routes-add.jinja
../../../../../debian/network/etc/network/if-up.d/routes-add.jinja